# TP 5/yt.service
[Unit]
Description="Service pour faire tourner yt-v2.sh"

[Service]
ExecStart=/srv/yt/yt-v2.sh
User=yt

[Install]
WantedBy=multi-user.target
